# com.hyperspace.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.hyperspace.bridge</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/env</string>
		<string>node</string>
		<string>/Users/USERNAME/.hyperspace-bridge/hyperspace-bridge.js</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/USERNAME/Library/Logs/hyperspace-bridge.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/USERNAME/Library/Logs/hyperspace-bridge.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/USERNAME/.hyperspace-bridge</string>
</dict>
</plist>
